# platform-ubuntu-2404/telegraf.ini
[program:telegraf]
command=/bin/sh -c "sleep 60 && /usr/bin/telegraf --config /etc/telegraf/telegraf.conf"
user=telegraf
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/telegraf.log
stderr_logfile=/var/log/supervisor/telegraf_error.log
startsecs=10
